# estudiantes/vitale-florencia/tp/mod-generador/etc/cron.d/generador
* * * * * /opt/src/cron-generador.sh
#Ejecutar el script cada tantos minutos, horas, dias, meses, dias de la semana
