# estudiantes/solange_moreyra/tp/generador/etc/cron.d/generator
 1 * * * * * /opt/src/cron-generator 10
